# launchd/com.personalos.cadence-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.personalos.cadence-weekly</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-lc</string>
		<string>/Users/sulaimanrahman/projects/personal-os/mcp-servers/cadence-draft.sh weekly</string>
	</array>

	<!-- Sundays 17:00 — drafts the just-finished ISO week, ready for the Sunday /week review. -->
	<!-- Opt-in: NOT auto-loaded. Enable on the ONE designated host only:                       -->
	<!--   cp launchd/com.personalos.cadence-weekly.plist ~/Library/LaunchAgents/                -->
	<!--   launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.personalos.cadence-weekly.plist -->
	<key>StartCalendarInterval</key>
	<array>
		<dict><key>Weekday</key><integer>0</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
	</array>

	<key>StandardOutPath</key>
	<string>/tmp/personalos-cadence-weekly.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/personalos-cadence-weekly.log</string>
</dict>
</plist>
