# launchd/com.personalos.daily-brief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.personalos.daily-brief</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-lc</string>
		<string>/Users/sulaimanrahman/projects/personal-os/mcp-servers/daily-brief.sh</string>
	</array>

	<!-- Weekday mornings at 6:00am, so the DRAFT is ready before you wake. -->
	<!-- Eager path. The lazy fallback is the SessionStart-ensures hook, which kicks -->
	<!-- the brief on first open if this didn't run (laptop asleep, weekend, etc.). -->
	<!-- The single mkdir-lock in daily-brief.sh guarantees the two never double-run. -->
	<!-- Weekend mornings are intentionally skipped (no eager spend); the hook covers -->
	<!-- them lazily if you open Claude. Weekday = launchd 1..5 (Mon-Fri). -->
	<key>StartCalendarInterval</key>
	<array>
		<dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
		<dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
		<dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
		<dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
		<dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
	</array>

	<key>StandardOutPath</key>
	<string>/tmp/personalos-daily-brief.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/personalos-daily-brief.log</string>
</dict>
</plist>
